# contrib/systemd/mail.target
[Unit]
Description=Mail Target

[Install]
WantedBy=default.target
