# contrib/systemd/offlineimap-oneshot@.timer
[Unit]
Description=Offlineimap Query Timer for account %i

[Timer]
OnBootSec=1m
OnUnitInactiveSec=15m

[Install]
WantedBy=default.target
