# contrib/systemd/offlineimap.timer
[Unit]
Description=Offlineimap Query Timer

[Timer]
OnUnitInactiveSec=15m
Unit=offlineimap.service

[Install]
WantedBy=mail.target
