# contrib/systemd/offlineimap@.timer
[Unit]
Description=Offlineimap Query Timer for account %i

[Timer]
OnUnitInactiveSec=15m
Unit=offlineimap@%i.service

[Install]
WantedBy=mail.target
