# contrib/systemd/offlineimap.service
[Unit]
Description=Offlineimap Service

[Service]
Type=oneshot
ExecStart=/usr/bin/offlineimap -o -s -u quiet

[Install]
WantedBy=mail.target
