# contrib/systemd/offlineimap.service
[Unit]
Description=Offlineimap Service

[Service]
ExecStart=/usr/bin/offlineimap -u syslog

[Install]
WantedBy=mail.target
